# systemd/brotherizer-api.service
[Unit]
Description=Brotherizer API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/genesis/brotherizer
EnvironmentFile=%h/genesis/brotherizer/.runtime/brotherizer.env
Environment=BROTHERIZER_HOST=127.0.0.1
Environment=BROTHERIZER_PORT=8787
ExecStart=/usr/bin/python3 /home/neodark/genesis/brotherizer/api/brotherizer_api.py
Restart=on-failure
RestartSec=2
StandardOutput=append:%h/genesis/brotherizer/.runtime/brotherizer_api.log
StandardError=append:%h/genesis/brotherizer/.runtime/brotherizer_api.log

[Install]
WantedBy=default.target
